# config/interject-scan.service
[Unit]
Description=Interject ambient discovery scanner
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/projects/Interverse/plugins/interject
ExecStart=/root/.local/bin/uv run --directory /root/projects/Interverse/plugins/interject interject-scan
TimeoutStartSec=600
Environment=HOME=/root
Environment=PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=EXA_API_KEY=eba9629f-75e9-467c-8912-a86b3ea8d678
